# homelab/company-lab/examples/company-backend.service.example
[Unit]
Description=Company Lab Backend API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=company-backend
Group=company-backend
WorkingDirectory=/opt/company-backend
EnvironmentFile=/etc/company-backend/company-backend.env
ExecStart=/usr/bin/python3 /opt/company-backend/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
